# osect_sensor/Infrastructure/edge_cron/work/crontab
RANDOM_TIME=$(($RANDOM%60))
* * * * * sleep $RANDOM_TIME; /opt/ot_tools/ot_cron.sh > /dev/null 2>&1
0 3 * * * /opt/ot_tools/complete_to_archives.sh > /dev/null 2>&1