# osect_sensor/Infrastructure/edge_tcpdump/work/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/bin/python3 /opt/sc_tcpdump/manage.py tcpdump_control

